# dist/openwrt/utils/lightsd/files/lightsd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 Louis Opter <kalessin@kalessin.fr>
#
# This is free software, licensed under the GNU General Public License v2.
# See /LICENSE for more information.

START=60

USE_PROCD=1

lightsd=/usr/bin/lightsd
rundir=$($lightsd --rundir)

start_service() {
	procd_open_instance
	procd_set_param command $lightsd -S -v warning -t -u lightsd -s "$rundir/socket" -c "$rundir/pipe"
	procd_set_param respawn
	procd_close_instance
}
