# dist/lightsd.service
[Unit]
Description=LIFX WiFi smart bulbs control service
After=network.target

[Service]
ExecStart=/usr/bin/lightsd -t -v warning -s %t/lightsd/socket -c %t/lightsd/pipe
User=lightsd
Group=lightsd
RuntimeDirectory=lightsd
Restart=on-failure

[Install]
WantedBy=multi-user.target
