# dist/lightsd.service
[Unit]
Description=LIFX WiFi smart bulbs control service
After=network.target

[Service]
ExecStart=/usr/bin/lightsd -t -v warning -s %t/lightsd/socket -c %t/lightsd/pipe
User=lightsd
Group=lightsd
RuntimeDirectory=lightsd
Restart=on-failure
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
# Debian Stretch and after:
# RestrictNamespaces=yes

[Install]
WantedBy=multi-user.target
